# debug/droplet-3/my-debug-script.service
[Unit]
Description=Debugging service configured as part of the initialization script

[Service]
ExecStart=/usr/bin/my-debug-script.sh
