# services/my-shutdown-script.service
[Unit]
Description=A special shutdown script that must be triggered manually

[Service]
ExecStart=/usr/bin/my-shutdown-script.sh